# scripts/com.cog.articles.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cog.articles</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/darrenarmitage/Library/Mobile Documents/iCloud~md~obsidian/Documents/Second Brain/scripts/run_articles.sh</string>
    </array>

    <!-- Run nightly at 23:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/cog-articles.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cog-articles-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
