# scripts/com.cog.vault-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cog.vault-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/darrenarmitage/Library/Mobile Documents/iCloud~md~obsidian/Documents/Second Brain/scripts/run_vault_review.sh</string>
    </array>

    <!-- Run daily at 06:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the machine was asleep at 06:00, run when it next wakes -->
    <key>StartOnMount</key>
    <false/>

    <!-- Log file for troubleshooting -->
    <key>StandardOutPath</key>
    <string>/tmp/cog-vault-review.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cog-vault-review-error.log</string>

    <!-- Do not run immediately on load — wait for 06:00 -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
